# deploy/ec2/callanalyzer.service
# Observatory QA — systemd service unit
# Install to: /etc/systemd/system/callanalyzer.service
#
# Commands:
#   sudo systemctl enable callanalyzer   # Start on boot
#   sudo systemctl start callanalyzer    # Start now
#   sudo systemctl status callanalyzer   # Check status
#   sudo journalctl -u callanalyzer -f   # Tail logs

[Unit]
Description=Observatory QA - AI-Powered Call Quality Analysis
After=network.target
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
User=callanalyzer
Group=callanalyzer
WorkingDirectory=/opt/callanalyzer
ExecStart=/usr/bin/node dist/index.js
Restart=always
RestartSec=5

# Environment
Environment=NODE_ENV=production
Environment=PORT=5000
EnvironmentFile=/opt/callanalyzer/.env

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/callanalyzer
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Resource limits
LimitNOFILE=65535
MemoryMax=1G

# Logging (goes to journald → accessible via journalctl)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=callanalyzer

[Install]
WantedBy=multi-user.target
